# systemd/system/sovereign-adapters-api.service
[Unit]
Description=Sovereign adapters read-only API + D-11 dashboard host (M060 R10109-R10111)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/adapters-api.py
# Read-only observation surface: joins the model catalog's LoRA-adapter entries
# (M046 LoRA Foundry) to the promotion registry (per-adapter status + MS041
# triple-gate + eval gain + history), surfacing the data model the D-11 cockpit
# dashboard renders. There is NO mutation surface — adapter promote/demote/
# rollback are MS003-signed CLI verbs gated by the MS041 triple-gate.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/adapters-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default — operator-§1g exposure decision is theirs. To expose
# beyond loopback, drop /etc/systemd/system/sovereign-adapters-api.service.d/bind.conf:
#   [Service]
#   Environment=ADAPTERS_API_BIND=0.0.0.0
Environment=ADAPTERS_API_BIND=127.0.0.1
Environment=ADAPTERS_API_PORT=8107

# R171 defense-in-depth (Layer-1 lint enforced). Read-only HTTP surface;
# reads the model catalog + adapter registry + appends the metrics textfile.
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
