# systemd/system/sovereign-gatewayd.service
[Unit]
Description=Sovereign gateway daemon — M048 Module 4 provider-inversion gateway over the cortex engine
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/crates/sovereign-gatewayd/src/lib.rs
# The first persistent runnable service: promotes the sovereign-cortex engine
# into a long-lived daemon behind the M048 gateway contract. Doctrine:
# "client -> Sovereign Gateway -> local/cloud/model router". Forces every
# request local (never-cloud-spill) and serves the manifest bind paths over
# HTTP: GET /health /manifest /admin/ledger /metrics; POST /v1/messages.
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sovereign-gatewayd --http
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default — the operator-§1g exposure decision is theirs. To expose
# beyond loopback, drop /etc/systemd/system/sovereign-gatewayd.service.d/bind.conf:
#   [Service]
#   Environment=SOVEREIGN_GATEWAY_ADDR=0.0.0.0:8787
Environment=SOVEREIGN_GATEWAY_ADDR=127.0.0.1:8787
# Durable memory: the daemon resumes from + periodically snapshots the learning
# Cortex here (atomic writes), so recall survives a restart instead of seeding
# cold each boot. StateDirectory grants the single writable path under strict.
Environment=SOVEREIGN_GATEWAY_MEMORY=/var/lib/sovereign-os/memory/cortex.json
StateDirectory=sovereign-os/memory

# R171 defense-in-depth (Layer-1 lint enforced). ProtectSystem=strict keeps the
# whole filesystem read-only EXCEPT the StateDirectory above (the memory
# snapshot) — the only write surface the daemon has.
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
