# systemd/system/sovereign-profile-mirror-api.service
[Unit]
Description=Sovereign profile-mirror read-only API + D-02 dashboard host (M060 R10063-R10068)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/profile-mirror-api.py
# CROSS-REPO MIRROR — renders selfdef's MS007 active-profile mirror READ-ONLY.
# The authoritative profile-authority state lives in selfdef (MS040 six-profile
# matrix + MS039 L0..L6 + Ring 0..4). There is NO mutation surface — profile
# switch is `sovereign profile set` / selfdefctl + MS003 verbs on the IPS side
# only (MS043 R10212).
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/profile-mirror-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default. To expose beyond loopback, drop
# /etc/systemd/system/sovereign-profile-mirror-api.service.d/bind.conf:
#   [Service]
#   Environment=PROFILE_MIRROR_API_BIND=0.0.0.0
Environment=PROFILE_MIRROR_API_BIND=127.0.0.1
Environment=PROFILE_MIRROR_API_PORT=8117

# R171 defense-in-depth (Layer-1 lint enforced). Read-only HTTP surface; reads
# the selfdef-published mirror artifact + appends the metrics textfile.
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
