# systemd/system/sovereign-super-model-api.service
[Unit]
Description=Sovereign super-model-manifest read-only API + D-19 dashboard host (M060 R10124-R10125)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/super-model-api.py
# Read-only observation surface (sovereign-os-native): computes the super-model
# manifest LIVE from the repo — git HEAD version + the milestone catalog
# (M001..M080 ids/titles/R-row counts) + config/super-model-manifest.toml
# editorial overlay (M053 11 build-phases + per-milestone family/status). There
# is NO mutation surface — the manifest is computed, not set.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/super-model-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default. To expose beyond loopback, drop
# /etc/systemd/system/sovereign-super-model-api.service.d/bind.conf:
#   [Service]
#   Environment=SUPER_MODEL_API_BIND=0.0.0.0
Environment=SUPER_MODEL_API_BIND=127.0.0.1
Environment=SUPER_MODEL_API_PORT=8119

# R171 defense-in-depth (Layer-1 lint enforced). Read-only HTTP surface; reads
# the milestone catalog + manifest config + git + appends the metrics textfile.
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
