# systemd/system/sovereign-pulse.service
[Unit]
Description=sovereign-os Pulse — ternary CPU inference (bitnet.cpp on CCD 0)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/docs/sdd/011-inference-backend-stack.md
After=network-online.target
Wants=network-online.target
# Pulse doesn't depend on GPU drivers; runs alongside Logic/Oracle

[Service]
Type=simple
EnvironmentFile=-/etc/sovereign-os/inference-pulse.env
ExecStart=/opt/sovereign-os/scripts/inference/start-pulse.sh
Restart=on-failure
RestartSec=10s
# CCD 0 cores pinned by taskset inside the start script.
# Memory limit: bitnet.cpp models < 10GB; cap to leave room for ZFS ARC + GPU mappings
MemoryHigh=12G
MemoryMax=16G
# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/sovereign-os /var/lib/sovereign-os
ReadOnlyPaths=/mnt/vault/models
PrivateTmp=true

[Install]
WantedBy=multi-user.target
