# systemd/system/sovereign-zfs-scrub.service
[Unit]
Description=ZFS scrub (weekly)
After=zfs.target
Requires=zfs.target

[Service]
Type=oneshot
ExecStart=/opt/sovereign-os/scripts/hooks/recurrent/zfs-scrub.sh
# Scrub kicks off async; the script returns once started. Don't hold
# systemd waiting for the full scrub (can take hours).
TimeoutStartSec=120s
# Log to journal
StandardOutput=journal
StandardError=journal
# Defense-in-depth — zpool/zfs commands only need /sys + /dev/zfs;
# script writes metrics under textfile_collector
ProtectSystem=strict
ReadWritePaths=/var/lib/node_exporter/textfile_collector /var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true
