# systemd/system/sovereign-auth-tier-api.service
[Unit]
Description=Sovereign auth-tier read-only REST API (R501, E11.M7++)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/auth-tier-api.py
# This API surface MIRRORS the data the CLI exposes via
# `sovereign-osctl auth-tier <verb>`. Mutation verb `set` stays
# CLI-only (operator §17 sacrosanct sovereignty boundary).
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/auth-tier-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default — operator-§1g exposure decision is theirs.
# To expose beyond loopback, drop an override at
# /etc/systemd/system/sovereign-auth-tier-api.service.d/bind.conf:
#   [Service]
#   Environment=AUTH_TIER_API_BIND=0.0.0.0
Environment=AUTH_TIER_API_BIND=127.0.0.1
Environment=AUTH_TIER_API_PORT=8091

# Defense-in-depth — this is a read-only HTTP surface; the daemon
# only needs read on auth-tier.py + the config file and append on the
# metrics textfile collector dir.
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
# Loopback default — capability lockdown matches.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
