# systemd/system/sovereign-network-vlan.service
[Unit]
Description=sovereign-os network VLAN configuration (one-shot, first boot)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/post-install/network-vlan-config.sh
ConditionFirstBoot=yes
# SAIN-01 bare-metal hardware setup — skip cleanly on a VM/emulator (no target
# GPU / dual-NIC / ZFS pool to touch). Runs on the real machine (not virtualized).
ConditionVirtualization=no
After=sovereign-friction-audit.service systemd-networkd.service
Wants=systemd-networkd.service
PartOf=sovereign-firstboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sovereign-os/active-profile.env
ExecStart=/opt/sovereign-os/scripts/hooks/post-install/network-vlan-config.sh
StandardOutput=journal
StandardError=journal
# Defense-in-depth — writes systemd-networkd config; needs CAP_NET_ADMIN
# but no unrelated filesystem write
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/etc/systemd/network /var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=sovereign-firstboot.target
