# systemd/system/sovereign-feature-test-lab-api.service
[Unit]
Description=Sovereign Feature Test Lab read-only REST API + panel (live per-feature testing)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/feature-test-lab-api.py
# Read-only observability: serves the panel and runs per-feature self-tests
# (sovereign-feature-selftest) on demand. A self-test is a read-only computation
# (no host mutation); do_POST fail-closes 405.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/feature-test-lab-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default — operator-§1g exposure decision is theirs.
# To expose beyond loopback, drop an override at
# /etc/systemd/system/sovereign-feature-test-lab-api.service.d/bind.conf:
#   [Service]
#   Environment=FEATURE_TEST_LAB_API_BIND=0.0.0.0
Environment=FEATURE_TEST_LAB_API_BIND=127.0.0.1
Environment=FEATURE_TEST_LAB_API_PORT=8135

# Defense-in-depth — read-only HTTP surface; the daemon reads the webapp +
# config and shells the sovereign-feature-selftest binary (a pure compute self
# test — no host mutation). R171 defense-in-depth (Layer 1 lint enforced).
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
