# systemd/system/sovereign-firstboot.service
[Unit]
Description=sovereign-os first-boot completion marker
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/docs/src/lifecycle/post-install.md
ConditionFirstBoot=yes
After=sovereign-friction-audit.service sovereign-vfio-bind.service sovereign-network-vlan.service sovereign-tetragon-policy-load.service sovereign-zfs-arc-clamp.service sovereign-nvidia-driver-bind.service sovereign-workstation-shell-setup.service
Requires=sovereign-friction-audit.service sovereign-network-vlan.service sovereign-tetragon-policy-load.service sovereign-zfs-arc-clamp.service
PartOf=sovereign-firstboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /var/lib/sovereign-os
ExecStart=/bin/sh -c 'date -u --iso-8601=seconds > /var/lib/sovereign-os/first-boot-complete'
ExecStartPost=/bin/sh -c 'for f in .vfio-bind-needs-reboot .nvidia-bind-needs-reboot; do [ -f "/var/lib/sovereign-os/$f" ] && echo "NOTICE: reboot required for GPU binding to take effect ($f)" > /dev/console; done; true'
StandardOutput=journal
StandardError=journal
# Defense-in-depth — completion marker only writes /var/lib/sovereign-os
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=sovereign-firstboot.target
