# systemd/system/sovereign-tetragon-policy-load.service
[Unit]
Description=sovereign-os Tetragon TracingPolicy load (one-shot, first boot)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/post-install/tetragon-policy-load.sh
ConditionFirstBoot=yes
After=tetragon.service zfs.target
Wants=tetragon.service zfs.target
PartOf=sovereign-firstboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sovereign-os/active-profile.env
ExecStart=/opt/sovereign-os/scripts/hooks/post-install/tetragon-policy-load.sh
StandardOutput=journal
StandardError=journal
# Defense-in-depth — copies policy files to /etc/tetragon; appends to
# /mnt/vault/context security_audit.log on drift.
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/etc/tetragon /mnt/vault/context /var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=sovereign-firstboot.target
