# systemd/system/sovereign-security-update-check.service
[Unit]
Description=sovereign-os daily security-update availability check
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/recurrent/security-update-check.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/sovereign-os/scripts/hooks/recurrent/security-update-check.sh
TimeoutStartSec=120s
StandardOutput=journal
StandardError=journal
# Defense-in-depth — runs 'apt update' (needs /var/lib/apt) + emits metrics
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
# StateDirectory/LogsDirectory CREATE their /var dirs; '-' tolerates
# absent host paths — a bare ReadWritePaths entry refuses to spawn the
# unit (226/NAMESPACE) on hosts missing the dir (first host bring-up
# 2026-06-12).
StateDirectory=sovereign-os
ReadWritePaths=/var/lib/apt /var/cache/apt -/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
