# systemd/system/sovereign-nvidia-power-limit.service
[Unit]
Description=sovereign-os NVIDIA per-GPU power caps (every boot)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/post-install/nvidia-power-limit.sh
# SAIN-01 bare-metal — skip cleanly on a VM/emulator (no GPU to cap).
ConditionVirtualization=no
# Apply after the driver is up. This runs EVERY boot (NOT ConditionFirstBoot) —
# nvidia-smi -pl does not persist across reboots, so the caps must be re-applied.
After=sovereign-nvidia-driver-bind.service
Wants=sovereign-nvidia-driver-bind.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sovereign-os/active-profile.env
ExecStart=/opt/sovereign-os/scripts/hooks/post-install/nvidia-power-limit.sh
StandardOutput=journal
StandardError=journal
# Only queries + sets nvidia-smi power limits (talks to /dev/nvidia*, writes the
# Layer-B metric). No system-file writes → ProtectSystem=strict + a narrow RWP.
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector /var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
