# systemd/system/sovereign-wattage-sample.service
[Unit]
Description=sovereign-os per-minute PSU wattage sampler (R258 / SDD-026 Z-18)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/recurrent/wattage-sample.sh
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/sovereign-os/scripts/hooks/recurrent/wattage-sample.sh
TimeoutStartSec=30s
StandardOutput=journal
StandardError=journal
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
# StateDirectory/LogsDirectory CREATE their /var dirs; '-' tolerates
# absent host paths — a bare ReadWritePaths entry refuses to spawn the
# unit (226/NAMESPACE) on hosts missing the dir (first host bring-up
# 2026-06-12).
StateDirectory=sovereign-os
ReadWritePaths=-/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
