# systemd/system/sovereign-oracle-core.service
[Unit]
Description=sovereign-os Oracle Core — vLLM + DFlash on Blackwell (host-resident)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/docs/sdd/011-inference-backend-stack.md
After=network-online.target nvidia-persistenced.service tetragon.service
Wants=network-online.target
Requires=tetragon.service

[Service]
Type=simple
EnvironmentFile=-/etc/sovereign-os/inference-oracle-core.env
ExecStart=/opt/sovereign-os/scripts/inference/start-oracle-core.sh
Restart=on-failure
RestartSec=20s
# Oracle Core loads can be slow (large weights + DFlash draft); generous startup
TimeoutStartSec=600s
# Memory pressure: vLLM uses GPU VRAM, but host buffers also; leave headroom
MemoryHigh=32G
# Security — defense in depth (R160 hardening pass)
# Same caveats as logic-engine: this service runs vLLM on Blackwell and
# needs CUDA / NVIDIA driver access; we cannot fully lock down /sys or
# kernel modules. ProtectHome=read-only allows runtime libs to read
# user config without writable paths.
NoNewPrivileges=true
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectKernelModules=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=false  # vLLM internals + podman compatibility
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
MemoryDenyWriteExecute=false  # vLLM/Triton + DFlash speculative-decode draft JIT-compile kernels; W^X would break codegen
ReadWritePaths=/var/log/sovereign-os /var/lib/sovereign-os
ReadOnlyPaths=/mnt/vault/models
PrivateTmp=true

[Install]
WantedBy=multi-user.target
