# systemd/system/sovereign-doc-coverage-api.service
[Unit]
Description=Sovereign doc-coverage read-only REST API + webapp (R527, E5++)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/operator/doc-coverage-api.py
# Read-only inspection API + webapp for the §1g doc-coverage instrument
# (R454 scanner). Mirrors the data the CLI exposes via
# `sovereign-osctl doc-coverage {kinds,modules,coverage,scan,gaps,selfdef}
# --json` and the R526 MCP tools (doc-coverage-kinds / -modules /
# -coverage).
#
# Doc-coverage has NO mutation verbs at any surface — docs ARE the
# source of truth; this daemon scans the tree, it does not author.
# Operator §17 sovereignty boundary preserved.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/sovereign-os/scripts/operator/doc-coverage-api.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Loopback-by-default — operator-§1g exposure decision is theirs.
# To expose beyond loopback, drop an override at
# /etc/systemd/system/sovereign-doc-coverage-api.service.d/bind.conf:
#   [Service]
#   Environment=DOC_COVERAGE_API_BIND=0.0.0.0
Environment=DOC_COVERAGE_API_BIND=127.0.0.1
Environment=DOC_COVERAGE_API_PORT=8099

# Defense-in-depth — read-only HTTP surface that walks the source tree
# per /scan, /coverage, /gaps request, plus the R471 selfdef
# DocManifest discovery on /selfdef. No mutation of any sovereignty
# boundary.
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector
# Source tree is read-only here — the scanner walks it to enumerate
# doc-surface presence but never mutates it. ProtectSystem=strict
# makes everything outside ReadWritePaths read-only by default.
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
# Loopback default — capability lockdown matches.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
