# systemd/system/sovereign-nvidia-driver-bind.service
[Unit]
Description=sovereign-os NVIDIA driver bind + nouveau blacklist (one-shot, first boot)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/post-install/nvidia-driver-bind.sh
ConditionFirstBoot=yes
After=sovereign-friction-audit.service
PartOf=sovereign-firstboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sovereign-os/active-profile.env
ExecStart=/opt/sovereign-os/scripts/hooks/post-install/nvidia-driver-bind.sh
StandardOutput=journal
StandardError=journal
# Defense-in-depth — writes /etc/modprobe.d for nouveau blacklist;
# reads /sys/bus/pci. Can't ProtectSystem=strict because kmod operations
# need /sys + module-loading capability. ReadWritePaths narrowed.
ProtectSystem=full
ReadWritePaths=/etc/modprobe.d /var/lib/sovereign-os
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=sovereign-firstboot.target
