# systemd/system/sovereign-warp-setup.service
[Unit]
Description=sovereign-os NVIDIA Warp setup — install warp-lang (one-shot, first boot)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/post-install/warp-setup.sh
ConditionFirstBoot=yes
# No ConditionVirtualization gate: warp-lang installs + runs on CPU too, so this
# is safe (and useful) on a VM/emulator, unlike the GPU-binding first-boot units.
After=sovereign-nvidia-driver-bind.service
PartOf=sovereign-firstboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sovereign-os/active-profile.env
ExecStart=/opt/sovereign-os/scripts/hooks/post-install/warp-setup.sh
StandardOutput=journal
StandardError=journal
# R171 defense-in-depth. pip writes the warp-lang wheel into the system
# python's dist-packages under /usr/local/lib; PrivateTmp gives pip its own
# build/cache dir; the metrics dir receives the Layer B result sample.
ProtectSystem=strict
ReadWritePaths=/usr/local/lib /usr/local/bin -/var/lib/node_exporter/textfile_collector /var/lib/sovereign-os /var/log
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=sovereign-firstboot.target
