# systemd/system/sovereign-logic-engine.service
[Unit]
Description=sovereign-os Logic Engine — vLLM/llama.cpp on RTX 3090 (VFIO sandbox)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/docs/sdd/011-inference-backend-stack.md
After=network-online.target podman.service nvidia-persistenced.service
Wants=network-online.target
# Logic Engine needs the 3090 VFIO binding active; tetragon perimeter loaded first
After=tetragon.service
Requires=tetragon.service

[Service]
Type=simple
EnvironmentFile=-/etc/sovereign-os/inference-logic-engine.env
ExecStart=/opt/sovereign-os/scripts/inference/start-logic-engine.sh
ExecStop=/usr/bin/podman stop vllm-logic_engine
Restart=on-failure
RestartSec=15s
TimeoutStartSec=180s
# Security — defense in depth (R160 hardening pass)
# Note: Logic Engine uses GPU via podman; we DO NOT add ProtectHome
# since podman may need access to /root/.config in some setups, and
# DO NOT lock down NVIDIA driver paths (vLLM needs CUDA introspection).
NoNewPrivileges=true
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectKernelModules=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
ProtectHome=read-only
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=false  # podman creates namespaces internally
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
MemoryDenyWriteExecute=false  # vLLM/Triton JIT-compile CUDA kernels at request time; W^X would break codegen
ReadWritePaths=/var/log/sovereign-os /var/lib/sovereign-os
ReadOnlyPaths=/mnt/vault/models
PrivateTmp=true

[Install]
WantedBy=multi-user.target
