# systemd/system/sovereign-models-sync.service
[Unit]
Description=Verify resident model catalog integrity
After=zfs.target

[Service]
Type=oneshot
ExecStart=/opt/sovereign-os/scripts/hooks/recurrent/model-catalog-sync.sh
TimeoutStartSec=300s
StandardOutput=journal
StandardError=journal
# Defense-in-depth — reads /mnt/vault/models; may write per-model
# verification ledgers + metrics output
# R171 defense-in-depth (Layer 1 lint enforced)
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
# StateDirectory/LogsDirectory CREATE their /var dirs; '-' tolerates
# absent host paths — a bare ReadWritePaths entry refuses to spawn the
# unit (226/NAMESPACE) on hosts missing the dir (first host bring-up
# 2026-06-12).
StateDirectory=sovereign-os
ReadWritePaths=-/mnt/vault/models -/var/lib/node_exporter/textfile_collector
NoNewPrivileges=true
PrivateTmp=true
