# systemd/system/sovereign-ghostproxy-verify.timer
[Unit]
Description=Weekly root-modules endpoint-envelope drift verify (SDD-046)
Documentation=https://github.com/cyberpunk042/sovereign-os/blob/main/scripts/hooks/recurrent/root-modules-verify.sh

[Timer]
OnCalendar=Sun 06:00:00
RandomizedDelaySec=1800
Persistent=true
Unit=sovereign-ghostproxy-verify.service

[Install]
WantedBy=timers.target
