# docs/observability/dashboards/sovereign-os-selfdef-systemd-units.json
{
  "_comment": "sovereign-os selfdef systemd-units-health dashboard. Visualizes the selfdef-side selfdef_systemd_units_* gauges shipped via the selfdef-systemd-units-textfile.{service,timer} pair (selfdef commit 7121c72). Pairs with the 4 alert rules. Silent-unit-failure detection — when a selfdef-* timer or service silently fails, this dashboard surfaces it before downstream observability degrades.",
  "annotations": {"list": []},
  "editable": true,
  "fiscalYearStartMonth": 0,
  "graphTooltip": 0,
  "id": null,
  "links": [
    {"title": "operator runbook (sovereign-os)", "url": "/d/sovereign-os-selfdef-systemd-units#runbook", "type": "link", "icon": "doc"},
    {"title": "producer source", "url": "https://github.com/cyberpunk042/selfdef/blob/main/packaging/scripts/selfdef-systemd-units-textfile.sh", "type": "link", "icon": "external link"},
    {"title": "deployment guide runbooks", "url": "https://github.com/cyberpunk042/sovereign-os/blob/main/docs/operator/m060-deployment-guide.md#selfdefsystemdunitsextfileemitfailed-critical", "type": "link", "icon": "external link"}
  ],
  "panels": [
    {
      "type": "stat",
      "title": "selfdef units total",
      "description": "Total selfdef-prefixed systemd units. Healthy floor is 10+. Yellow at 8 (matches SelfdefSystemdUnitsCountLow).",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [{"expr": "selfdef_systemd_units_total", "refId": "A"}],
      "fieldConfig": {"defaults": {"unit": "none", "thresholds": {"mode": "absolute", "steps": [{"color": "red", "value": 0}, {"color": "yellow", "value": 8}, {"color": "green", "value": 10}]}}},
      "gridPos": {"h": 4, "w": 6, "x": 0, "y": 0}
    },
    {
      "type": "stat",
      "title": "active",
      "description": "selfdef units currently active. Green = healthy operating state.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [{"expr": "selfdef_systemd_units_active", "refId": "A"}],
      "fieldConfig": {"defaults": {"unit": "none", "thresholds": {"mode": "absolute", "steps": [{"color": "red", "value": 0}, {"color": "green", "value": 1}]}}},
      "gridPos": {"h": 4, "w": 6, "x": 6, "y": 0}
    },
    {
      "type": "stat",
      "title": "failed",
      "description": "selfdef units in failed state. Red at 1 = page-worthy (matches SelfdefSystemdUnitFailed).",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [{"expr": "selfdef_systemd_units_failed", "refId": "A"}],
      "fieldConfig": {"defaults": {"unit": "none", "thresholds": {"mode": "absolute", "steps": [{"color": "green", "value": 0}, {"color": "red", "value": 1}]}}},
      "gridPos": {"h": 4, "w": 6, "x": 12, "y": 0}
    },
    {
      "type": "stat",
      "title": "observer emit-failed",
      "description": "Honest-offline sentinel.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [{"expr": "selfdef_systemd_units_textfile_emit_failed", "refId": "A"}],
      "fieldConfig": {"defaults": {"unit": "none", "min": 0, "max": 1, "mappings": [{"type": "value", "options": {"0": {"text": "OK"}, "1": {"text": "FAILED"}}}], "thresholds": {"mode": "absolute", "steps": [{"color": "green", "value": 0}, {"color": "red", "value": 1}]}}},
      "gridPos": {"h": 4, "w": 6, "x": 18, "y": 0}
    },
    {
      "type": "timeseries",
      "title": "per-state counts over time",
      "description": "All 5 ActiveState gauges side-by-side. A failed-count spike = silent unit failure. Inactive count drift = operator-disabled units.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [
        {"expr": "selfdef_systemd_units_active", "refId": "A", "legendFormat": "active"},
        {"expr": "selfdef_systemd_units_inactive", "refId": "B", "legendFormat": "inactive"},
        {"expr": "selfdef_systemd_units_failed", "refId": "C", "legendFormat": "failed"},
        {"expr": "selfdef_systemd_units_activating", "refId": "D", "legendFormat": "activating"},
        {"expr": "selfdef_systemd_units_other", "refId": "E", "legendFormat": "other"}
      ],
      "fieldConfig": {"defaults": {"unit": "none", "custom": {"drawStyle": "line", "lineInterpolation": "stepAfter", "lineWidth": 2, "fillOpacity": 30, "stacking": {"mode": "normal", "group": "A"}}}},
      "gridPos": {"h": 8, "w": 12, "x": 0, "y": 4}
    },
    {
      "type": "timeseries",
      "title": "total units over time (with 8-floor threshold)",
      "description": "Total selfdef unit count. Red threshold at 8 matches SelfdefSystemdUnitsCountLow alert floor.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [
        {"expr": "selfdef_systemd_units_total", "refId": "A", "legendFormat": "{{instance}}"}
      ],
      "fieldConfig": {"defaults": {"unit": "none", "custom": {"drawStyle": "line", "lineInterpolation": "stepAfter", "lineWidth": 2, "fillOpacity": 30}, "thresholds": {"mode": "absolute", "steps": [{"color": "red", "value": 0}, {"color": "yellow", "value": 8}, {"color": "green", "value": 10}]}}},
      "gridPos": {"h": 8, "w": 12, "x": 12, "y": 4}
    },
    {
      "type": "timeseries",
      "title": "failed-unit count over time (page-worthy spikes)",
      "description": "Dedicated panel for the failed-unit count. ANY spike above 0 is page-worthy and signals downstream observer degradation.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [
        {"expr": "selfdef_systemd_units_failed", "refId": "A", "legendFormat": "{{instance}}"}
      ],
      "fieldConfig": {"defaults": {"unit": "none", "min": 0, "custom": {"drawStyle": "line", "lineInterpolation": "stepAfter", "lineWidth": 2, "fillOpacity": 50}, "thresholds": {"mode": "absolute", "steps": [{"color": "green", "value": 0}, {"color": "red", "value": 1}]}}},
      "gridPos": {"h": 7, "w": 12, "x": 0, "y": 12}
    },
    {
      "type": "timeseries",
      "title": "observer-freshness — emit-failed sentinel",
      "description": "Red bands = wedged-wrapper windows.",
      "datasource": {"type": "prometheus", "uid": "prometheus"},
      "targets": [
        {"expr": "selfdef_systemd_units_textfile_emit_failed", "refId": "A", "legendFormat": "{{instance}}"}
      ],
      "fieldConfig": {"defaults": {"unit": "none", "min": 0, "max": 1, "custom": {"drawStyle": "line", "lineInterpolation": "stepAfter", "lineWidth": 2, "fillOpacity": 50}, "thresholds": {"mode": "absolute", "steps": [{"color": "green", "value": 0}, {"color": "red", "value": 1}]}}},
      "gridPos": {"h": 7, "w": 12, "x": 12, "y": 12}
    }
  ],
  "refresh": "30s",
  "schemaVersion": 39,
  "style": "dark",
  "tags": ["sovereign-os", "selfdef", "systemd-units", "IPS-spine", "observability"],
  "templating": {"list": []},
  "time": {"from": "now-1h", "to": "now"},
  "timepicker": {},
  "timezone": "",
  "title": "sovereign-os — selfdef systemd units (silent-failure detection)",
  "uid": "sovereign-os-selfdef-systemd-units",
  "version": 1,
  "weekStart": ""
}
